# .config/systemd/user/BAT.timer
[Unit]
Description=BAT

[Timer]
Unit=BAT.service
Persistent=false
OnCalendar=*-*-* *:*:00
#AccuracySec=1s

[Install]
WantedBy=timers.target
