# .config/systemd/user/BAT.service
[Unit]
Description=BAT

[Service]
Type=oneshot
ExecStart=%h/.config/systemd/user/BAT.sh
WorkingDirectory=%h

[Install]
WantedBy=default.target
