# .config/systemd/user/BAT0.timer
[Unit]
Description=BAT0

[Timer]
Unit=BAT0.service
Persistent=false
OnCalendar=*-*-* *:*:01
#AccuracySec=1s

[Install]
WantedBy=timers.target
