# .config/systemd/user/BAT0.service
[Unit]
Description=BAT0

[Service]
Type=simple
ExecStart=%h/.dotfiles/linux/BAT0.sh
WorkingDirectory=%h

[Install]
WantedBy=default.target
