# .config/systemd/user/borg.timer
[Unit]
Description=borg

[Timer]
Unit=borg.service
Persistent=false
OnCalendar=*-*-* *:05:00

[Install]
WantedBy=timers.target
