# .config/systemd/user/firefox.service
[Unit]
Description=firefox

[Service]
Type=simple
ExecStart=/bin/sh -c ./backup.sh
WorkingDirectory=%h

[Install]
WantedBy=default.target
