# .config/systemd/user/xscreensaver.service
[Unit]
Description=Lock X session
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/xscreensaver-command -lock

[Install]
WantedBy=sleep.target
