# .config/autostart/autocutsel.desktop
[Desktop Entry]
Type=Application
Name=autocutsel
Exec=autocutsel -s PRIMARY
TryExec=autocutsel
Terminal=false
