# .config/autostart/devilspie.desktop
[Desktop Entry]
Type=Application
Name=devilspie
Exec=devilspie2
TryExec=devilspie2
Terminal=false
