# .config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=dunst
Exec=dunst
TryExec=dunst
Terminal=false
