# .config/autostart/easystroke.desktop
[Desktop Entry]
Type=Application
Name=easystroke
Exec=easystroke
TryExec=easystroke
Terminal=false
