# .config/autostart/gxkb.desktop
[Desktop Entry]
Name=gxkb
Icon=gxkb
Exec=gxkb
Terminal=false
Type=Application
