# .config/autostart/gxkb.desktop
[Desktop Entry]
Type=Application
Name=gxkb
Exec=gxkb
Terminal=false
