# .config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=nm-applet
Exec=nm-applet
Terminal=false
