# .config/autostart/opensnitch-ui.desktop
[Desktop Entry]
Type=Application
Name=opensnitch-ui
Exec=opensnitch-ui
TryExec=opensnitch-ui
Terminal=false
