# .config/autostart/pasystray.desktop
[Desktop Entry]
Name=pasystray
Icon=pasystray
Exec=pasystray
Type=Application
