# .config/autostart/pasystray.desktop
[Desktop Entry]
Type=Application
Name=pasystray
Exec=pasystray
TryExec=pasystray
Terminal=false
