# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Icon=redshift
Exec=redshift-gtk -l 50:8
Terminal=false
