# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Icon=redshift
Exec=redshift-gtk -l 50:8 -b 0.75:0.75
TryExec=redshift-gtk
Terminal=false
