# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Type=Application
Name=redshift
Icon=redshift
Exec=redshift-gtk -l 50:8 -b 0.7:0.7
TryExec=redshift-gtk
Terminal=false
