# .config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Name=sxhkd
Icon=sxhkd
Exec=sh -c "$DOTFIELS/bin/sxhkd"
TryExec=redshift-gtk
Terminal=false
