# .config/autostart/uxterm.desktop
[Desktop Entry]
Name=uxterm
Exec=sh -c 'uxterm'
Type=Application
