# .config/autostart/uxterm.desktop
[Desktop Entry]
Type=Application
Name=uxterm
Exec=sh -c 'uxterm'
