# .config/autostart/xss-lock.desktop
[Desktop Entry]
Type=Application
Name=xss-lock
Exec=xss-lock -l -- lock.sh
TryExec=xss-lock
Terminal=false
