# profile.d/00-path.sh
#!/bin/sh

my_path_push_front() {
    case ":$PATH:" in
        *:"$1":*) ;;
        ::) PATH="$1" ;;
        *) PATH="$1:$PATH" ;;
    esac
}
