# profile.d/50-mu3e.sh
#

export ROOTSYS="$HOME/.local/root.cern.ch"

export MIDASSYS="$HOME/mu3e/midas"
export MIDAS_EXPT_NAME=mu3e
#export MIDAS_SERVER_HOST=127.0.0.1
export MIDAS_EXPTAB="$MIDASSYS/exptab"
