# profile.d/50-mpv.sh
#!/bin/sh

export MPV_SOCKET="/run/user/$(id -u)/mpv.sock"
