# profile.d/90-cmake.sh
#

if command -v ninja &> /dev/null ; then
    export CMAKE_GENERATOR=Ninja
fi
