# profile.d/90-debug.sh
#!/bin/sh

export DEBUGINFOD_URLS="https://debuginfod.archlinux.org/"
