# profile.d/90-python.sh
#!/bin/sh

export MY_PYTHON_VENV="${XDG_CACHE_HOME}/venv"
