# profile.d/90-pico.sh
#!/bin/sh

export PICO_SDK_PATH="$XDG_CACHE_HOME/pico-sdk"
export OPENOCD_PICOPROBE="$XDG_CACHE_HOME/openocd-picoprobe"
