# profile.d/90-rust.sh
#

# <https://doc.rust-lang.org/cargo/reference/environment-variables.html>

# location of where to place all generated artifacts, relative to the current working directory
export CARGO_TARGET_DIR=".cache"
