# profile.d/90-v4p.sh
#!/bin/sh

export V4P_MASTER_CONFIG="$XDG_CONFIG_HOME/v4p/config.v4p"
export V4P_LICENSE_FILE="$XDG_CONFIG_HOME/v4p/license.v4p"
